# conf/etc/s6/services/0000-crond/run
#!/bin/sh

service-handler "CRON" "crond -f" "/tmp/container.init.end"