# conf/etc/s6/services/0000-crond/run
#!/bin/bash

service-handler "CRON" "cron -f" "/tmp/container.init.end"